# deploy/djafatt.service
[Unit]
Description=djafatt gunicorn
After=network.target postgresql.service redis-server.service

[Service]
User=www
Group=www
WorkingDirectory=/www/wwwroot/fatt.betabi.it
EnvironmentFile=/www/wwwroot/fatt.betabi.it/.env
ExecStart=/www/wwwroot/fatt.betabi.it/venv/bin/gunicorn djafatt.wsgi:application \
    --bind 127.0.0.1:8200 \
    --workers 3 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile -
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
